# scripts/templates/gtbi-nightly-update.timer
[Unit]
Description=GTBI Nightly Auto-Update Timer
Documentation=https://github.com/jonbackhaus/gtbi

[Timer]
# Run daily at 4am local time
OnCalendar=*-*-* 04:00:00

# If machine was off at 4am, run on next boot
Persistent=true

# Spread load across machines (0-5 min jitter)
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
