# scripts/templates/gtbi-upgrade-resume.service
[Unit]
Description=GTBI Ubuntu Upgrade Resume Service
Documentation=https://github.com/jonbackhaus/gtbi

# Wait for network to be fully online before running
# This is critical because do-release-upgrade needs apt to work
After=network-online.target
Wants=network-online.target

# Only run if the resume script exists (upgrade in progress)
ConditionPathExists=/var/lib/gtbi/upgrade_resume.sh

# Don't run in rescue mode or during recovery
ConditionPathExists=!/run/systemd/maintenance

[Service]
Type=oneshot

# Run the resume script
ExecStart=/bin/bash /var/lib/gtbi/upgrade_resume.sh

# Give plenty of time for upgrades (2 hours max)
TimeoutStartSec=7200

# Don't restart on failure - let user investigate
Restart=no
RemainAfterExit=no

# Logging to journal
StandardOutput=journal+console
StandardError=journal+console

# Run as root (required for apt and reboot)
User=root
Group=root

[Install]
WantedBy=multi-user.target
